# AdapterMonitor/crontab
* * * * * /root/check_c.sh 
* * * * * /root/check_r.sh

* * * * * /root/adapter_check.sh