# AdapterMonitor/crontab
* * * * * /root/check_c.sh 
* * * * * /root/check_r.sh